# infrastructure/vm-cron/crontab
# BuildAtlas VM Cron Jobs
# Installed via: crontab infrastructure/vm-cron/crontab
# All times are UTC
#
# runner.sh handles: .env sourcing, flock locking, timeout, logging, Slack alerts

SHELL=/bin/bash
PATH=/opt/buildatlas/venv/bin:/usr/local/bin:/usr/bin:/bin
MAILTO=""

# === Infrastructure Keep-Alive (every 15 min) ===
# Checks: PostgreSQL → AKS → API health → Frontend
*/15 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh keep-alive 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/keep-alive.sh

# === News Ingest (hourly at :15) ===
15 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh news-ingest 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/news-ingest.sh

# === Crawl Frontier Worker (every 30 min) ===
0,30 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh crawl-frontier 25 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/crawl-frontier.sh

# === Topic Research (hourly at :40, 25 min after news-ingest) ===
# Processes hot topic research queue: web search + LLM synthesis
40 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh research-topics 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/research-topics.sh

# === News Digest (hourly at :45, timezone-aware) ===
# Sends to subscribers whose local time is 08:xx. Each run only delivers to
# the timezone batch that matches; duplicates are prevented by delivery table.
45 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh news-digest 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/news-digest.sh

# === Weekly Brief (Monday 06:00 UTC) ===
0 6 * * 1 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh weekly-brief 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/weekly-brief.sh

# === Monthly Brief (1st of month 06:00 UTC) ===
0 6 1 * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh monthly-brief 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/monthly-brief.sh

# === Embedding Backfill (hourly at :25, after news-ingest) ===
# Embeds unembedded clusters; at 03:xx UTC also populates related_cluster_ids
25 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh embed-backfill 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/embed-backfill.sh

# === Signal Aggregation (every 4 hours at :30, after news-ingest) ===
# Aggregates events into signals, scores, and updates lifecycle
30 */4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh signal-aggregate 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/signal-aggregate.sh

# === Dealbook Living Brief (daily 04:00 UTC / 07:00 Istanbul) ===
# Generates fresh brief snapshots for each region
0 4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh dealbook-brief 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/dealbook-brief.sh

# === Digest QA Preview (hourly at :50) ===
# Sends merged global+turkey digest to QA email for iterative review
50 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh digest-qa 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/digest-qa.sh

# === Infrastructure Health Report (every 4 hours at :45) ===
# Posts a concise all-systems status to Slack (6 reports/day)
45 0,4,8,12,16,20 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh health-report 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/health-report.sh

# === Daily Observability Report (daily 09:00 UTC) ===
# Pipeline quality SLO report (latency, refresh, linking, embedding coverage)
0 9 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh daily-observability 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/daily-observability.sh

# === Slack Daily Summary (daily 14:00 UTC) ===
0 14 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh slack-summary 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/slack-summary.sh

# === Slack Commit Notify (every 2 min) ===
*/2 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh slack-commit-notify 5 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/slack-commit-notify.sh

# === Sync Blob Data (every 30 min, all days) ===
0,30 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh sync-data 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/sync-data.sh

# === Frontend Deploy (triggered by sync-data or manually) ===
# Not scheduled — called by sync-data.sh after a push, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh frontend-deploy 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/frontend-deploy.sh

# === Backend Deploy (triggered by code-update or manually) ===
# Not scheduled — called by deploy.sh when apps/api/** changes, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh backend-deploy 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/backend-deploy.sh

# === Code Update (every 15 min, staggered) ===
# Runs git pull + conditional deploys. Staggered to avoid colliding with :00/:15/:30/:45 jobs.
7,22,37,52 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh code-update 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/deploy.sh

# === VM Health Monitoring (every 5 min) ===
*/5 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/monitoring/heartbeat.sh >> /var/log/buildatlas/heartbeat.log 2>&1
