# infrastructure/vm-cron/crontab
# BuildAtlas VM Cron Jobs
# Installed via: crontab infrastructure/vm-cron/crontab
# All times are UTC
#
# runner.sh handles: .env sourcing, flock locking, timeout, logging, Slack alerts

SHELL=/bin/bash
PATH=/opt/buildatlas/venv/bin:/usr/local/bin:/usr/bin:/bin
MAILTO=""

# === Infrastructure Keep-Alive (every 15 min) ===
# Checks: PostgreSQL → AKS → API health → Frontend
*/15 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh keep-alive 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/keep-alive.sh

# === News Ingest (hourly at :15) ===
15 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh news-ingest 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/news-ingest.sh

# === Crawl Frontier Worker (every 30 min) ===
0,30 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh crawl-frontier 25 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/crawl-frontier.sh

# === News Digest (daily 13:10 UTC) ===
10 13 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh news-digest 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/news-digest.sh

# === Slack Daily Summary (daily 14:00 UTC) ===
0 14 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh slack-summary 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/slack-summary.sh

# === Slack Commit Notify (every 2 min) ===
*/2 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh slack-commit-notify 5 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/slack-commit-notify.sh

# === Sync Blob Data (every 30 min, all days) ===
# Note: job is cheap when no changes are detected (it exits early after a --check).
0,30 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh sync-data 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/sync-data.sh

# === Frontend Deploy (triggered by sync-data or manually) ===
# Not scheduled — called by sync-data.sh after a push, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh frontend-deploy 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/frontend-deploy.sh

# === Backend Deploy (triggered by code-update or manually) ===
# Not scheduled — called by deploy.sh when apps/api/** changes, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh backend-deploy 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/backend-deploy.sh

# === Code Update (every 15 min, staggered) ===
# Runs git pull + conditional deploys. Staggered to avoid colliding with :00/:15/:30/:45 jobs.
7,22,37,52 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh code-update 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/deploy.sh

# === VM Health Monitoring (every 5 min) ===
*/5 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/monitoring/heartbeat.sh >> /var/log/buildatlas/heartbeat.log 2>&1
