# infrastructure/vm-cron/crontab
# BuildAtlas VM Cron Jobs
# Installed via: crontab infrastructure/vm-cron/crontab
# All times are UTC
#
# runner.sh handles: .env sourcing, flock locking, timeout, logging, Slack alerts

SHELL=/bin/bash
PATH=/opt/buildatlas/venv/bin:/usr/local/bin:/usr/bin:/bin
MAILTO=""

# === Infrastructure Keep-Alive (every 15 min) ===
# Checks: PostgreSQL → AKS → API health → Frontend
*/15 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh keep-alive 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/keep-alive.sh

# === News Ingest (hourly at :15) ===
15 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh news-ingest 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/news-ingest.sh

# === X Trend Ingest (hourly at :28) ===
# Pulls recent X/Twitter trend signals and runs clustering via ingest-x-trends.
28 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh x-trends 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/x-trends.sh

# === Event Processor (every 15 min, staggered) ===
# Processes startup_events and enqueues eligible deep-research jobs.
5,20,35,50 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh event-processor 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/event-processor.sh

# === Deep Research Consumer (every 15 min, staggered) ===
# Consumes deep_research_queue with strict budget caps.
12,27,42,57 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh deep-research 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/deep-research.sh

# === Investor Onboarding Consumer (every 30 min at :08,:38) ===
# Consumes investor_onboarding_queue with conservative budget caps.
8,38 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh investor-onboarding 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/investor-onboarding.sh

# === Onboarding Trace Alerts (every 2 min) ===
# Dispatches actionable onboarding/deep-research trace events to Slack.
*/2 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh onboarding-alerts 5 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/onboarding-alerts.sh

# === Crawl Frontier Worker (every 30 min) ===
0,30 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh crawl-frontier 40 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/crawl-frontier.sh

# === Topic Research (hourly at :40, 25 min after news-ingest) ===
# Processes hot topic research queue: web search + LLM synthesis
40 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh research-topics 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/research-topics.sh

# === News Digest (hourly at :45, timezone-aware) ===
# Sends to subscribers whose local time is 08:xx. Each run only delivers to
# the timezone batch that matches; duplicates are prevented by delivery table.
45 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh news-digest 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/news-digest.sh

# === X Post Candidate Generation (every 4h at :35) ===
35 */4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh x-post-generate 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/x-post-generate.sh

# === X Post Publish (hourly at :55) ===
55 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh x-post-publish 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/x-post-publish.sh

# === X Post Metrics Sync (every 6h at :20) ===
20 */6 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh x-post-metrics 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/x-post-metrics.sh

# === Weekly Brief (Monday 06:00 UTC) ===
0 6 * * 1 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh weekly-brief 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/weekly-brief.sh

# === Weekly Watchlist Digest (Monday 06:35 UTC) ===
# Groups recent alerts into digest threads. Best-effort LLM summary; falls back to heuristic summary if LLM unavailable.
35 6 * * 1 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh generate-weekly-digest 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/generate-weekly-digest.sh

# === Monthly Brief (1st of month 06:00 UTC) ===
0 6 1 * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh monthly-brief 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/monthly-brief.sh

# === Embedding Backfill (hourly at :25, after news-ingest) ===
# Embeds unembedded clusters; at 03:xx UTC also populates related_cluster_ids
25 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh embed-backfill 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/embed-backfill.sh

# === Signal Aggregation (every 4 hours at :30, after news-ingest) ===
# Aggregates events into signals, scores, and updates lifecycle
30 */4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh signal-aggregate 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/signal-aggregate.sh

# === Delta Events + Watchlist Alerts (every 4 hours, staggered after signal-aggregate) ===
# Generates delta_events (movers feed), then materializes user_alerts for watchlists/subscriptions.
38 */4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh delta-generate 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/delta-generate.sh
48 */4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh generate-alerts 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/generate-alerts.sh

# === Deep Dive Generation (daily 05:15 UTC) ===
# Computes occurrences and generates signal deep dives (moves + synthesis)
15 5 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh deep-dive-generate 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/deep-dive-generate.sh

# === Deep Dive Catchup (every 4 hours at :58 UTC) ===
# Backfills missing deep dives so all non-decaying signals have coverage content.
58 */4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh deep-dive-catchup 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/deep-dive-catchup.sh

# === Dealbook Living Brief (daily 04:00 UTC / 07:00 Istanbul) ===
# Generates fresh brief snapshots for each region
0 4 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh dealbook-brief 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/dealbook-brief.sh

# === Neighbors & Benchmarks (weekly Thu 07:00 UTC) ===
0 7 * * 4 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh neighbors-benchmarks 60 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/neighbors-benchmarks.sh

# === Extended Benchmarks (monthly 2nd 04:00 UTC) ===
0 4 2 * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh compute-benchmarks 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/compute-benchmarks.sh

# === Investor DNA (monthly 2nd 05:00 UTC) ===
0 5 2 * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh compute-investor-dna 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/compute-investor-dna.sh

# === Digest QA Preview (every 3 hours at :50) ===
# Sends merged global+turkey digest to QA email for iterative review
50 */3 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh digest-qa 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/digest-qa.sh

# === Product Canary (every 30 min) ===
# Verifies brief schema + Investor DNA + deep dives; posts Slack only on state changes/reminders.
17,47 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh product-canary 5 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/product-canary.sh

# === Infrastructure Health Report (every 4 hours at :45) ===
# Posts a concise all-systems status to Slack (6 reports/day)
45 0,4,8,12,16,20 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh health-report 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/health-report.sh

# === Daily Observability Report (daily 09:00 UTC) ===
# Pipeline quality SLO report (latency, refresh, linking, embedding coverage)
0 9 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh daily-observability 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/daily-observability.sh

# === EOD Onboarding/Graph/News Report (daily 20:00 UTC) ===
# Posts a Slack report of startups/investors added, graph edge updates, and news<>startup linking activity.
0 20 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh onboarding-eod-report 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/onboarding-eod-report.sh

# === Slack + Metrics Email Summary (every 3 hours) ===
0 */3 * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh slack-summary 10 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/slack-summary.sh

# === Slack Commit Notify (every 2 min) ===
*/2 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh slack-commit-notify 5 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/slack-commit-notify.sh

# === Release Reconciler (every 5 min) ===
# Compares desired main SHA with live frontend/backend build SHAs and posts Slack on drift.
*/5 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh release-reconciler 5 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/release-reconciler.sh

# === Sync Blob Data (every 30 min, all days) ===
0,30 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh sync-data 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/sync-data.sh

# === Frontend Deploy (triggered by sync-data or manually) ===
# Not scheduled — called by sync-data.sh after a push, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh frontend-deploy 20 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/frontend-deploy.sh

# === Backend Deploy (triggered by code-update or manually) ===
# Not scheduled — called by deploy.sh when apps/api/** changes, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh backend-deploy 15 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/backend-deploy.sh

# === Functions Deploy (triggered by code-update or manually) ===
# Not scheduled — called by deploy.sh when infrastructure/azure-functions/** or
# packages/analysis/** changes, or run manually:
#   /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh functions-deploy 30 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/jobs/functions-deploy.sh

# === Code Update (every 15 min, staggered) ===
# Runs git pull + conditional deploys. Staggered to avoid colliding with :00/:15/:30/:45 jobs.
7,22,37,52 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/lib/runner.sh code-update 45 /opt/buildatlas/startup-analysis/infrastructure/vm-cron/deploy.sh

# === VM Health Monitoring (every 5 min) ===
*/5 * * * * /opt/buildatlas/startup-analysis/infrastructure/vm-cron/monitoring/heartbeat.sh >> /var/log/buildatlas/heartbeat.log 2>&1
